# com.wavetree.scrapbook-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wavetree.scrapbook-backup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>NEWS_HUB_DIR="${WAVETREE_NEWS_HUB_DIR:-${WAVETREE_WORKSPACE_ROOT:-$HOME/AI_WORKSPACE}/wave-tree-news-hub}"; export BACKUP_DIR="$NEWS_HUB_DIR/data/scrapbook"; export SYNC_DIR="$NEWS_HUB_DIR/data/scrapbook"; exec node "$NEWS_HUB_DIR/backup_server.js"</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/tmp</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/wavetree-backup-server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/wavetree-backup-server-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
