# com.wavetree.perplexity-auto.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
DEPRECATED (2026-02-21)
- LaunchAgent 기반 운영 중단: 외장 볼륨 경로 접근 시 Operation not permitted 발생 가능
- 현재 운영 기준: cron (06:50 run_perplexity_auto.sh, 07:00 run_daily_bridge.sh)
- 참고 문서: AUTOMATION_STATUS.md
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.wavetree.perplexity-auto</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>NEWS_HUB_DIR="${WAVETREE_NEWS_HUB_DIR:-${WAVETREE_WORKSPACE_ROOT:-$HOME/AI_WORKSPACE}/wave-tree-news-hub}"; exec /bin/bash "$NEWS_HUB_DIR/run_perplexity_auto.sh"</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/wavetree-perplexity-auto.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/wavetree-perplexity-auto-error.log</string>
	<key>WorkingDirectory</key>
	<string>/tmp</string>
</dict>
</plist>
