# com.wavetree.httpserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wavetree.httpserver</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>bash</string>
        <string>-c</string>
        <string>cd /Users/seunghoonoh/Desktop/wave-tree-news-hub && /usr/bin/python3 -m http.server 8000</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/wavetree-http.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/wavetree-http-error.log</string>
</dict>
</plist>
